# crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user  command
55 23     * * *   root    python /app/cronjob.py >> /var/log/cron.log 2>&1

